# etc/profullstack-pdf.service
[Unit]
Description=Document Generation Service
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/www/profullstack.com/pdf
ExecStart=/bin/bash %h/www/profullstack.com/pdf/bin/start.sh
Restart=on-failure
RestartSec=10

# Log to files instead of journal
StandardOutput=journal
StandardError=journal

# Environment
Environment=NODE_ENV=production
Environment=PATH=/usr/local/bin:/usr/bin:/bin:/usr/local/sbin:/usr/sbin:/sbin:%h/.local/share/pnpm:%h/.npm/pnpm/bin

# Hardening
ProtectSystem=full
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target