# etc/profullstack-pdf.service
[Unit]
Description=Document Generation Service
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/www/profullstack.com/pdf
ExecStart=/bin/zsh /home/ubuntu/www/profullstack.com/pdf/bin/start.sh
Restart=on-failure
RestartSec=10

# Log to files instead of journal
StandardOutput=append:/var/log/profullstack-pdf.log
StandardError=append:/var/log/profullstack-pdf.error.log

# Environment
Environment=NODE_ENV=production
Environment=PATH=/usr/local/bin:/usr/bin:/bin:/usr/local/sbin:/usr/sbin:/sbin:/home/ubuntu/.local/share/pnpm:/home/ubuntu/.npm/pnpm/bin
Environment=HOME=/home/ubuntu

# Hardening
ProtectSystem=full
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target